# s6/etc/s6-overlay/s6-rc.d/web/run
#!/command/with-contenv sh
frankenphp php-server --root /app/public --listen :8000
