# s6/etc/s6-overlay/s6-rc.d/cache/up
#!/command/with-contenv sh
php /app/bin/console cache:clear
