# s6/etc/s6-overlay/s6-rc.d/initialize/up
#!/command/with-contenv sh
frankenphp php-cli /app/bin/console app:initialize
