# s6/etc/s6-overlay/s6-rc.d/worker/run
#!/command/with-contenv sh
frankenphp php-cli /app/bin/console messenger:consume scheduler_default async --limit=10 --time-limit=600
