# s6/etc/s6-overlay/s6-rc.d/scheduler/run
#!/command/with-contenv sh
php /app/bin/console messenger:consume scheduler_default --time-limit=86400
